# deploy/launchd/local.pxx.gpu-node-1-vllm-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.pxx.gpu-node-1-vllm-tunnel</string>

    <!--
      macOS launchd analog of inference-node's coder-tunnel.service.
      Persistent SSH local-forward from the Mac Studio to the gpu-node-1
      audit-proxy (port 8003), which fronts the 2x RTX A4500 vLLM serving
      qwen2.5-coder-14b. The gpu-node-1 is SSH-only (office NAT forwards
      only port 22), so this tunnel is BOTH the transport and the security
      boundary -- the audit-proxy relies on SSH for access control.

      Result: http://127.0.0.1:8003 on the Studio -> gpu-node-1 vLLM.
      pxx reads DEFAULT_VLLM / PXX_VLLM_URL to find it.

      Install: see deploy/README.md. Replace /Users/you with your $HOME
      and confirm an `ssh gpu-node-1` alias exists in ~/.ssh/config.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-N</string>
        <string>-F</string>
        <string>/Users/you/.ssh/config</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=accept-new</string>
        <string>-L</string>
        <string>8003:127.0.0.1:8003</string>
        <string>gpu-node-1</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/you</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/you/Library/Logs/pxx-gpu-node-1-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/you/Library/Logs/pxx-gpu-node-1-tunnel.log</string>
</dict>
</plist>
